# deployment/systemd.service
# Servicio Systemd para gestion_mantenimiento
# Copia este archivo a /etc/systemd/system/gestion_mantenimiento.service
# Habilita: sudo systemctl enable gestion_mantenimiento
# Inicia: sudo systemctl start gestion_mantenimiento
# Reinicia: sudo systemctl restart gestion_mantenimiento

[Unit]
Description=Gestion Mantenimiento Django App
After=network.target

[Service]
User=tu-usuario  # Cambia por tu usuario del sistema
Group=www-data
WorkingDirectory=/ruta/a/tu/proyecto  # Cambia por la ruta absoluta al directorio del proyecto
Environment="PATH=/ruta/a/tu/venv/bin"  # Cambia por la ruta a tu virtualenv
ExecStart=/ruta/a/tu/venv/bin/gunicorn --workers 3 --bind unix:/ruta/a/tu/proyecto/gestion_mantenimiento.sock gestion_mantenimiento.wsgi:application
Restart=always

[Install]
WantedBy=multi-user.target